# supervisor.ini
[supervisord]
nodaemon=true

[program:flaskweb]
command=gunicorn -c /webapp/conf/gunicorn.conf app:app
directory=/webapp/web
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/webapp/log/flaskout.log
stderr_logfile=/webapp/log/flaskerr.log

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true